# door-proxy/digilog-door-proxy.service
[Unit]
Description=Digilog Door Proxy (WebSocket-to-Telnet bridge)
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/html/YOURDOMAIN-app/door-proxy
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
